# systemd/jetmon2.service
[Unit]
Description=Jetmon 2 — Jetpack Uptime Monitor
Documentation=https://github.com/Automattic/jetmon
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=jetmon
Group=jetmon
WorkingDirectory=/opt/jetmon2
ExecStart=/opt/jetmon2/jetmon2
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=35s
StartLimitIntervalSec=60s
StartLimitBurst=5

# Resource limits.
MemoryMax=512M
LimitNOFILE=65536
LimitNPROC=4096

# Hardening.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes

# RuntimeDirectory creates /run/jetmon2/ owned by the service user.
RuntimeDirectory=jetmon2
RuntimeDirectoryMode=0750

# Environment.
EnvironmentFile=-/opt/jetmon2/config/jetmon2.env
Environment=JETMON_CONFIG=/opt/jetmon2/config/config.json

# PID file for jetmon2 reload/drain subcommands.
PIDFile=/run/jetmon2/jetmon2.pid

[Install]
WantedBy=multi-user.target
