# systemd/jetmon-deliverer.service
[Unit]
Description=Jetmon Delivery Workers
Documentation=https://github.com/Automattic/jetmon
After=network.target mysql.service
Wants=network.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
User=jetmon
Group=jetmon
WorkingDirectory=/opt/jetmon2
ExecStartPre=/opt/jetmon2/bin/jetmon-deliverer validate-config --require-owner-match --require-api-disabled
ExecStart=/opt/jetmon2/bin/jetmon-deliverer
Restart=on-failure
RestartSec=5s
TimeoutStopSec=35s

# Resource limits.
MemoryMax=256M
LimitNOFILE=65536
LimitNPROC=4096

# Hardening.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes

# Environment.
EnvironmentFile=-/opt/jetmon2/config/jetmon2.env
Environment=JETMON_CONFIG=/opt/jetmon2/config/deliverer.json

[Install]
WantedBy=multi-user.target
